# com.blender-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.blender-mcp</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Blender.app/Contents/MacOS/Blender</string>
		<string>--background</string>
		<string>--python-expr</string>
		<string>import bpy; bpy.ops.wm.blender_mcp_server()</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/ryanstern/foundry/logs/blender-mcp.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/ryanstern/foundry/logs/blender-mcp.stderr.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/ryanstern/foundry</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
